# scripts/systemd/autolab-idle-eval.service
[Unit]
Description=Autolab idle evaluator loop (up to 14 CPUs)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/www
ExecStart=/www/scripts/autolab_idle_eval_loop.sh
Restart=always
RestartSec=15
Nice=10

[Install]
WantedBy=multi-user.target
